# simple-test-opsmanager-k8s/src/supervisor-mms.conf
[supervisord]
nodaemon=true

[program:on-on-k8s]
command=/mongodb-mms/mms-launcher.sh
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
