# simple-test-opsmanager-k8s/src/supervisor-mongod.conf
[supervisord]
nodaemon=true

[program:mongodb-opsmgr-appdb]
command=/mongodb-opsmgr/mongod-launcher.sh
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
